# backend/conversation-ai.service
[Unit]
Description=Conversation AI Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/conversation-ai
ExecStart=/usr/bin/npm start
Restart=on-failure
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
